# templates/puppetmaster-webhook.service.erb
[Unit]
Description=Puppetmaster Webhook

[Service]
Environment=APP_DIR=/var/lib/puppet
Environment=APP_EXEC=/usr/local/bin/puppetmaster-webhook
Environment=USER=puppet
Environment=BINDIP=0.0.0.0
Environment=BINDPORT=8141
Environment=TOKEN=a1b9c114cd650bc60c53
ExecStart=/sbin/start-stop-daemon --start --chuid "${USER}" --chdir "${APP_DIR}" --exec "${APP_EXEC}" -- ${BINDIP} ${BINDPORT} ${TOKEN}
#ExecStart=/sbin/start-stop-daemon --start --chuid puppet --chdir /var/lib/puppet --exec /usr/local/bin/puppetmaster-webhook -- 0.0.0.0 8141 a1b9c114cd650bc60c53
ExecStop=/bin/true

[Install]
WantedBy=multi-user.target
